# templates/sensor_rcscript.epp
#!/bin/ksh
#
# $OpenBSD$

daemon="/usr/local/sbin/<%= $agenttype %>_agent.tcl"
daemon_flags="-D -c <%= $configfile %>"
daemon_user="<%= $daemon_user %>"

. /etc/rc.d/rc.subr

pexp="/usr/local/bin/tclsh8.5 ${daemon}${daemon_flags:+ ${daemon_flags}}"

rc_cmd $1
